# config/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/sbin/php-fpm7.0

[program:nginx]
command=/usr/sbin/nginx